# units/master/kube-controller-manager.service
[Service]
ExecStart=/opt/bin/kube-controller-manager --master=127.0.0.1:8080

[X-Fleet]
MachineMetadata=master=true
